# system-report.service
[Unit]
Description=System Report Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/system_report.sh
